# deploy/supervisor/twint-telegram.conf
; Supervisor unit for the Telegram operator bot (long-polling).
; Install:
;   sudo apt install -y supervisor
;   sudo cp deploy/supervisor/twint-telegram.conf /etc/supervisor/conf.d/twint-ubc-telegram.conf
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start twint-ubc-telegram
; Adjust command path, directory and user to your deployment.

[program:twint-ubc-telegram]
process_name=%(program_name)s
command=/usr/bin/php /var/www/twint-ubc/artisan telegram:poll
directory=/var/www/twint-ubc
user=www-data
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=15
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/supervisor/twint-ubc-telegram.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=HOME="/var/www",USER="www-data"
