# deploy/supervisor/twint-reverb.conf
; Supervisor unit for Laravel Reverb (WebSocket server).
; Install:
;   sudo apt install -y supervisor
;   sudo cp deploy/supervisor/twint-reverb.conf /etc/supervisor/conf.d/twint-ubc-reverb.conf
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start twint-ubc-reverb
; Adjust command path, directory and user to your deployment.

[program:twint-ubc-reverb]
process_name=%(program_name)s
command=/usr/bin/php /var/www/twint-ubc/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/twint-ubc
user=www-data
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=15
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/supervisor/twint-ubc-reverb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=HOME="/var/www",USER="www-data"
